# .service
[Unit]
Description=Internet Speed Monitoring
After=network.target

[Service]
ExecStart=java -jar /home/pi/services/internet-speed-monitor.jar
EnvironmentFile=/etc/internet-speed-monitor/internet-speed-monitor.conf
Restart=always
User=pi

[Install]
WantedBy=multi-user.target